# config/einheit-ui.service
[Unit]
Description=einheit-ui web server (takt adapter)
After=network.target takt-service.service

[Service]
Type=simple
ExecStart=%h/dev/takt/build/takt-ui \
  --takt-url http://127.0.0.1:7433 \
  --editor
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=%h/dev/takt

[Install]
WantedBy=default.target
