# config/takt-service.service
[Unit]
Description=takt background service
After=network.target

[Service]
Type=simple
ExecStart=%h/dev/takt/.venv/bin/python3 %h/dev/takt/bin/takt_service.py
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.nvm/versions/node/v24.11.1/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=%h/dev/takt

[Install]
WantedBy=default.target
